# home/.config/autostart/localsend_app.desktop
[Desktop Entry]
Type=Application
Name=localsend_app
Comment=localsend_app startup script
Exec=/usr/lib/localsend/localsend
StartupNotify=false
Terminal=false
